# hybrid-backend/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:rust-server]
command=/app/rust-server
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
environment=PORT="8890",PYTHON_AI_URL="http://127.0.0.1:8891"

[program:python-ai]
command=python -m uvicorn ai_server:app --host 0.0.0.0 --port 8891
directory=/app/python-ai
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
environment=DEEPSEEK_API_KEY="sk-2fbc148564eb48c28576314d2fa89a30"
